# dist/amane-client.service
[Unit]
Description=amane multipath bonding tunnel (client)
Documentation=https://github.com/TKYcraft/amane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/amane client -c /etc/amane/client.toml
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
