# dist/dev.tkycraft.amane.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS client daemon: copy to /Library/LaunchDaemons/ and run
     sudo launchctl bootstrap system /Library/LaunchDaemons/dev.tkycraft.amane.plist -->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>dev.tkycraft.amane</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/amane</string>
		<string>client</string>
		<string>-c</string>
		<string>/etc/amane/client.toml</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>NetworkState</key>
		<true/>
	</dict>
	<key>StandardErrorPath</key>
	<string>/var/log/amane.log</string>
</dict>
</plist>
